# engine/provisioner/bifrost/playbooks/roles/create-libvirt-vms/templates/vbmcd.service.j2
[Unit]
Description=Virtual BMC daemon

[Service]
Type=simple
Restart=on-failure
ExecStart={{ engine_venv }}/bin/vbmcd --foreground
